# packaging/slothingd.service
# systemd user service for the Slothing decode daemon (25M ternary, ggml/libslothe).
# Install:  packaging/install-slothingd-service.sh
# The socket lands on $XDG_RUNTIME_DIR/slothingd.sock, the same path the
# fcitx5/IBus engines derive (see slothingdSocketPath() in eim.cpp).
#
# Prereqs: build the daemon once —
#   cmake -S engine/slothingd -B engine/slothingd/build_slothe -DCMAKE_BUILD_TYPE=Release
#   cmake --build engine/slothingd/build_slothe --target slothingd_slothe
# and fetch the model — packaging/fetch-model.sh
[Unit]
Description=Slothing zhuyin decode daemon (SlothLM-E-T 25M ternary)

[Service]
Type=simple
# REPO_DIR is substituted by the install script.
ExecStart=@REPO_DIR@/engine/slothingd/build_slothe/slothingd_slothe \
    -m @REPO_DIR@/model/slothe_t_25m/slothe-t-25m.gguf \
    -t @REPO_DIR@/model/phonetic_table.tsv \
    -v @REPO_DIR@/model/slothe_t_25m/syl_vocab.json \
    -j @REPO_DIR@/model/slothe_t_25m/char2id.json
Restart=on-failure
RestartSec=2
# tiny ternary model; keep the daemon lean and out of the way
MemoryHigh=512M
Nice=5

[Install]
WantedBy=default.target
